# service/xtrlock.service
[Unit]
Description=xtrlock

[Service]
ExecStart=/usr/bin/xtrlock
Restart=always

[Install]
WantedBy=default.target
